# misc/.config/systemd/user/tocavpn.socket
[Socket]
ListenFIFO=%t/tocavpn.stdin
Service=tocavpn.service

[Install]
WantedBy=default.target
