# _nodes/mac/etc/systemd/system/ntp-sync.service
[Unit]
Description=Sync NTP

[Service]
Type=simple
ExecStart=/usr/bin/ntpdate -v -s us.pool.ntp.org

[Install]
WantedBy=multi-user.target
