# _nodes/mac/etc/systemd/system/ntp-sync.timer
[Unit]
Description=Sync time from ntp every 3 hours

[Timer]
OnBootSec=1min
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
