# _arch/etc/systemd/system/autoborg@root.timer
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar=*-*-* 10:20
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
