# _arch/etc/systemd/system/autoglacier.timer
[Unit]
Description=Run rclone glacier backup

[Timer]
OnCalendar=*-*-* 10:45
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
