# _arch/etc/systemd/system/closetomb.service
[Unit]
Description=Close tomb before shutdown
Requires=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=true
RemainAfterExit=true
ExecStop=bash -c "tomb slam || true"

[Install]
WantedBy=multi-user.target
