# _arch/etc/backup/autousb/replicate_backup@.service
[Service]
Type=oneshot
ExecStart=/etc/backup/autousb/replicate_backup.sh %i
