# _arch/etc/borgmatic/replicate_backup.service
[Service]
Type=oneshot
ExecStart=/etc/borgmatic/replicate_backup.sh
