# _arch/etc/backup/autoglacier/autoglacier.timer
[Unit]
Description=Run mtglacier backup

[Timer]
OnCalendar=*-*-* 10:45
Persistent=true

[Install]
WantedBy=timers.target
