# _arch/etc/systemd/system/https_dns_proxy.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/https_dns_proxy -v -r "https://doh.getbetter.ro/getnsrecord" -p 53

