# _arch/etc/systemd/system/autoborg@root.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 9:20