# _arch/etc/systemd/system/autoborg@home.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 9:00